# bot.conf
[program:bot]
directory = /opt/bot
command = /opt/bot/tg-channel-bot
autostart = true
autorestart = true
startsecs = 5
startretries = 30
user = nobody